# spot-the-gnome.service
[Unit]
Description=Spot the GNOME
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/spot-the-gnome/spot-the-gnome
